# templates/minikube.service.j2
[Unit]
Description=Minikube server



[Service]
Type=forking
ExecStart=/usr/local/bin/minikube start
User={{ansible_user_id}}




[Install]
WantedBy=multi-user.target